# example/rootfs/etc/services.d/example1/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting MochaD..."

exec /usr/bin/example1.sh
